# supervisord.conf
[supervisord]
nodaemon=true


[program:zentao]
command=/opt/zbox/zbox start && tail -f /dev/null
